# deploy/systemd/butfun-devloop.service
[Unit]
Description=ButFun autonomous dev loop (one iteration)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/butfun
Environment=BUTFUN_REPO=/opt/butfun
Environment=HOME=/home/butfun
# headless Claude Code 與 cargo 都要靠 PATH 找;把使用者常見路徑帶進來。
Environment=PATH=/home/butfun/.local/bin:/home/butfun/.cargo/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/butfun/scripts/devloop.sh
User=butfun
TimeoutStartSec=1800
