# deploy/systemd/butfun-devloop.timer
[Unit]
Description=Run ButFun dev loop nearly continuously (1 min gap between iterations)

[Timer]
OnBootSec=30s
# 從「上一輪結束」算 60 秒後再觸發 → 不會把仍在跑的下一輪疊起來
OnUnitInactiveSec=60s
AccuracySec=10s

[Install]
WantedBy=timers.target
