# deploy/systemd/butfun-devloop.timer
[Unit]
Description=Run ButFun dev loop nearly continuously (10s gap)

[Timer]
OnBootSec=15s
OnUnitInactiveSec=10s
AccuracySec=2s

[Install]
WantedBy=timers.target
