# deploy/systemd/butfun-deploy.timer
# prod 自動上線排程：啟用後約 2 分跑第一次，之後每 15 分一次。
# deploy.sh 沒新 commit 會自動跳過，所以這只是「最多多常重啟一次」（重啟會短暫踢線上玩家）。
# 改頻率就動 OnUnitActiveSec，重跑 scripts/deploy-enable.sh 生效。
[Unit]
Description=ButFun prod 自動上線（每 15 分，沒新 commit 自動跳過）

[Timer]
OnActiveSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
